# worker/decktr-worker.service
# systemd unit for the Decktr worker on the GC VPS.
# Install: copy to /etc/systemd/system/, `systemctl daemon-reload && systemctl enable --now decktr-worker`.
# EnvironmentFile holds the secrets (chmod 600, never committed).
[Unit]
Description=Decktr translation worker (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/deploy/ai-deck-translator
EnvironmentFile=/home/deploy/decktr-worker.env
ExecStart=/home/deploy/ai-deck-translator/.venv/bin/uvicorn worker.app:app --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
